# test/slurm/supervisor_docker.conf
[program:docker]
user=root
autostart=%(ENV_DOCKER_AUTOSTART)s
command=/usr/bin/docker daemon 
